# web/services/nginx/run
#!/bin/bash

chown -R www-data:www-data /data/www
chown -R www-data:www-data /data/log

/usr/bin/nginx 2>&1
