# postgres/services/postgres/run
#!/bin/bash

PGDATA=${PGDATA:-/data/pgdata}

exec 2>&1 \
  su -c "/usr/bin/postgres -D $PGDATA" - postgres
